# services/api/daphne.socket
[Unit]
Description=Daphne Socket

[Socket]
ListenStream=/run/daphne/daphne.sock

[Install]
WantedBy=sockets.target
